# system_files/usr/lib/systemd/user/wvkbd-deskintl.service
[Unit]
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target

[Service]
ExecStart=/usr/bin/wvkbd-deskintl --hidden
Restart=on-failure
