# system_files/usr/lib/systemd/user/wvkbd-mobintl.service
[Unit]
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target

[Service]
ExecStart=/usr/bin/wvkbd-mobintl --hidden
Restart=on-failure
