# system_files/usr/lib/systemd/user/niri.service.d/99-conflict-with-gamescope.conf
[Unit]
Conflicts=gamescope-session-plus@steam.service